# data/templates/konstant-scheduler.conf.erb
description "konstant ci scheduler (<%= Konstant.config['data_dir'] %>)"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid app
setgid app

script
  cd <%= Konstant.config['data_dir'] %>
  /opt/rbenv/shims/konstant run -s
end script
